# src/etc/cron.d/ebam-plus
# ensure runtime directories exist for scripts
@reboot root mkdir -p /var/run/ebam

# retrieve new data
0 * * * * root /home/lar/ebam-plus/src/getbam >> /home/lar/ebam-plus/src/getbam.log

# populate RPi Monitor status page
# but don't collide with existing data script
1-59 * * * * root python /home/lar/ebam-plus/src/test_ebam_comms.py

